# backend_django/systemd/daphne.service
[Unit]
Description=CampusCore Daphne ASGI Service
After=network.target redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/campuscore/backend_django
Environment=DJANGO_SETTINGS_MODULE=hostelconnect.settings.base
ExecStart=/var/www/campuscore/backend_django/venv/bin/daphne -b 127.0.0.1 -p 8001 hostelconnect.asgi:application
Restart=always
RestartSec=5
TimeoutStartSec=60
TimeoutStopSec=30
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

